# recipe-editor.service
[Unit]
Description=Recipe Editor Flask Application
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/home/yourusername/recipe_editor
Environment=PATH=/home/yourusername/recipe_editor/venv/bin
Environment=FLASK_ENV=production
Environment=STORAGE_BACKEND=mysql
ExecStart=/home/yourusername/recipe_editor/venv/bin/gunicorn --config gunicorn.conf.py app_mysql:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=recipe-editor

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/yourusername/recipe_editor/logs
ReadWritePaths=/home/yourusername/recipe_editor/data

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
